# deploy/systemd/user/invera.target.tmpl
[Unit]
Description=Invera full stack
Wants=invera-frontend.service invera-backend.service cloudflared.service
After=network-online.target

[Install]
WantedBy=default.target
